# 01-Python/green_linux/powerapp/cli/units/powerapp-rapl-helper.service
[Unit]
Description=PowerApp RAPL sampling helper
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/powerapp-rapl-helper
Restart=on-failure

[Install]
WantedBy=multi-user.target
